# scripts/goxlr-utility.service
[Unit]
Description=Handle GoXLR Utility initialization

[Service]
Type=oneshot
ExecStart=/usr/bin/goxlr-daemon

[Install]
WantedBy=default.target
